# crontab
* * * * * root bash /run-diffengine.sh

